# sqrdao-bot.service
[Unit]
Description=sqrDAO Telegram Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/sqrDAO-AI-Agent
Environment=PATH=/home/ubuntu/sqrDAO-AI-Agent/venv/bin
ExecStart=/home/ubuntu/sqrDAO-AI-Agent/venv/bin/python bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target 